# infrastructure/services/local/supervisord/comment.conf
[program:comment_create]
command=php bin/console swarrot:consume:comment_create comment_create --requeue-on-error --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
stdout_logfile=syslog
stderr_logfile=syslog

[program:comment_update]
command=php bin/console swarrot:consume:comment_update comment_update --requeue-on-error --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog

[program:comment_delete]
command=php bin/console swarrot:consume:comment_delete comment_delete --requeue-on-error --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog
