# infrastructure/services/local/supervisord/comment.conf
[program:comment_create]
command=php bin/console swarrot:consume:comment_create comment_create --requeue-on-error
directory=/var/www
autostart=true
autorestart=true
startsecs=10
startretries=15
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=1
user=capco

stdout_logfile=syslog
stderr_logfile=syslog

[program:comment_update]
command=php bin/console swarrot:consume:comment_update comment_update --requeue-on-error
directory=/var/www
process_name=alert_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=15
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=1
user=capco

stdout_logfile=syslog
stderr_logfile=syslog

[program:comment_delete]
command=php bin/console swarrot:consume:comment_delete comment_delete --requeue-on-error
directory=/var/www
process_name=alert_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=15
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=1
user=capco

stdout_logfile=syslog
stderr_logfile=syslog
