# infrastructure/services/local/supervisord/comment.conf
[program:comment_create]
command=php bin/console swarrot:consume:comment_create comment_create --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:comment_update]
command=php bin/console swarrot:consume:comment_update comment_update --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:comment_delete]
command=php bin/console swarrot:consume:comment_delete comment_delete --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:proposal_analysis_comment_create]
command=php bin/console swarrot:consume:proposal_analysis_comment_create proposal_analysis_comment_create --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
