# infrastructure/services/local/supervisord/debate_article_crawler.conf
[program:debate_article_crawler]
command=php bin/console swarrot:consume:debate_article_crawler debate_article_crawler --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
