# infrastructure/services/local/supervisord/elasticsearch.conf
[program:elasticsearch_indexation]
command=php bin/console swarrot:consume:elasticsearch_indexation elasticsearch_indexation --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
