# infrastructure/services/local/supervisord/opinion.conf
[program:opinion_create]
command=php bin/console swarrot:consume:opinion_create opinion_create --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog

[program:opinion_update]
command=php bin/console swarrot:consume:opinion_update opinion_update --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog

[program:opinion_delete]
command=php bin/console swarrot:consume:opinion_delete opinion_delete --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
