# infrastructure/services/local/supervisord/organization.conf
[program:organization_member_invitation]
command=php bin/console swarrot:consume:organization_member_invitation organization_member_invitation --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0