# infrastructure/services/local/supervisord/participant.conf
[program:participant_confirmation_email]
command=php bin/console swarrot:consume:participant_confirmation_email participant_confirmation_email --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:participant_consent_internal_communication_email]
command=php bin/console swarrot:consume:participant_consent_internal_communication_email participant_consent_internal_communication_email --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
