# infrastructure/services/local/supervisord/proposal.conf
[program:proposal_create]
command=php bin/console swarrot:consume:proposal_create proposal_create --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog

[program:proposal_create_dl]
command=php bin/console swarrot:consume:proposal_create proposal_create_dl --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_create_retry_25]
command=php bin/console swarrot:consume:proposal_create proposal_create_retry_25 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_create_retry_125]
command=php bin/console swarrot:consume:proposal_create proposal_create_retry_125 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_create_retry_625]
command=php bin/console swarrot:consume:proposal_create proposal_create_retry_625 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_update]
command=php bin/console swarrot:consume:proposal_update proposal_update --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_update_dl]
command=php bin/console swarrot:consume:proposal_update proposal_update_dl --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_update_retry_25]
command=php bin/console swarrot:consume:proposal_update proposal_update_retry_25 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_update_retry_125]
command=php bin/console swarrot:consume:proposal_update proposal_update_retry_125 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_update_retry_625]
command=php bin/console swarrot:consume:proposal_update proposal_update_retry_625 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog

[program:proposal_delete]
command=php bin/console swarrot:consume:proposal_delete proposal_delete --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_delete_dl]
command=php bin/console swarrot:consume:proposal_delete proposal_delete_dl --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_delete_retry_25]
command=php bin/console swarrot:consume:proposal_delete proposal_delete_retry_25 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_delete_retry_125]
command=php bin/console swarrot:consume:proposal_delete proposal_delete_retry_125 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog


[program:proposal_delete_retry_625]
command=php bin/console swarrot:consume:proposal_delete proposal_delete_retry_625 --requeue-on-error --no-catch --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startretries=10
stdout_logfile=syslog
stderr_logfile=syslog
