# infrastructure/services/local/supervisord/questionnaire.conf
[program:questionnaire_reply]
command=php bin/console swarrot:consume:questionnaire_reply questionnaire_reply --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
