# infrastructure/services/local/supervisord/report.conf
[program:report]
command=sh -c "php bin/console swarrot:consume:report report --env=test"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
