# infrastructure/services/local/supervisord/sms.conf
[program:sms_credit_initial_credit]
command=php bin/console swarrot:consume:sms_credit_initial_credit sms_credit_initial_credit --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_refill_credit]
command=php bin/console swarrot:consume:sms_credit_refill_credit sms_credit_refill_credit --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_initial_order]
command=php bin/console swarrot:consume:sms_credit_initial_order sms_credit_initial_order --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_refill_order]
command=php bin/console swarrot:consume:sms_credit_refill_order sms_credit_refill_order --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_alert_consumed_credit]
command=php bin/console swarrot:consume:sms_credit_alert_consumed_credit sms_credit_alert_consumed_credit --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0