# infrastructure/services/local/supervisord/user.conf
[program:user_password]
command=php bin/console swarrot:consume:user_password user_password --env=test
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
