# infrastructure/services/remote/supervisord/argument.conf
[program:argument_create]
command=php bin/console swarrot:consume:argument_create argument_create --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog

[program:argument_update]
command=php bin/console swarrot:consume:argument_update argument_update --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog

[program:argument_delete]
command=php bin/console swarrot:consume:argument_delete argument_delete --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
