# infrastructure/services/remote/supervisord/clamd.conf
[program:clamd]
command=/bin/bash -c "mkdir -p /var/run/clamav /var/log/clamav && chown -R clamav:clamav /var/run/clamav /var/log/clamav && rm -f /var/run/clamav/clamd.pid && exec /usr/sbin/clamd --foreground"
autostart=false
autorestart=true
user=root
startsecs=10
stdout_logfile=/var/log/clamav/clamd-stdout.log
stderr_logfile=/var/log/clamav/clamd-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
