# infrastructure/services/remote/supervisord/comment.conf
[program:comment_create]
command=php bin/console swarrot:consume:comment_create comment_create --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:comment_update]
command=php bin/console swarrot:consume:comment_update comment_update --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:comment_delete]
command=php bin/console swarrot:consume:comment_delete comment_delete --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
