# infrastructure/services/remote/supervisord/elasticsearch.conf
[program:elasticsearch_indexation]
command=php bin/console swarrot:consume:elasticsearch_indexation elasticsearch_indexation --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
