# infrastructure/services/remote/supervisord/event.conf
[program:event_create]
command=sh -c "php bin/console swarrot:consume:event_create event_create --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:event_update]
command=sh -c "php bin/console swarrot:consume:event_update event_update --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:event_delete]
command=sh -c "php bin/console swarrot:consume:event_delete event_delete --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
