# infrastructure/services/remote/supervisord/event.conf
[program:event_create]
command=php bin/console swarrot:consume:event_create event_create --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:event_update]
command=php bin/console swarrot:consume:event_update event_update --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:event_delete]
command=php bin/console swarrot:consume:event_delete event_delete --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
