# infrastructure/services/remote/supervisord/mediator.conf
[program:mediator_participate_email]
command=sh -c "php bin/console swarrot:consume:mediator_participate_email mediator_participate_email --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
