# infrastructure/services/remote/supervisord/messenger.conf
[program:messenger-consume]
command=sh -c "php bin/console messenger:consume async --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0