# infrastructure/services/remote/supervisord/participant.conf
[program:participant_confirmation_email]
command=sh -c "php bin/console swarrot:consume:participant_confirmation_email participant_confirmation_email --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:participant_consent_internal_communication_email]
command=sh -c "php bin/console swarrot:consume:participant_consent_internal_communication_email participant_consent_internal_communication_email --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
