# infrastructure/services/remote/supervisord/proposal.conf
[program:proposal_create]
command=php bin/console swarrot:consume:proposal_create proposal_create --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:proposal_update]
command=php bin/console swarrot:consume:proposal_update proposal_update --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:proposal_delete]
command=php bin/console swarrot:consume:proposal_delete proposal_delete --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:proposal_update_status]
command=php bin/console swarrot:consume:proposal_update_status proposal_update_status --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
