# infrastructure/services/remote/supervisord/proposal.conf
[program:proposal_create]
command=php bin/console swarrot:consume:proposal_create proposal_create --requeue-on-error --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
stdout_logfile=syslog
stderr_logfile=syslog

[program:proposal_update]
command=php bin/console swarrot:consume:proposal_update proposal_update --requeue-on-error --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
stdout_logfile=syslog
stderr_logfile=syslog

[program:proposal_delete]
command=php bin/console swarrot:consume:proposal_delete proposal_delete --requeue-on-error --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
stdout_logfile=syslog
stderr_logfile=syslog
