# infrastructure/services/remote/supervisord/questionnaire.conf
[program:questionnaire_reply]
command=php bin/console swarrot:consume:questionnaire_reply questionnaire_reply --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
