# infrastructure/services/remote/supervisord/sms.conf
[program:sms_credit_initial_credit]
command=sh -c "php bin/console swarrot:consume:sms_credit_initial_credit sms_credit_initial_credit --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_refill_credit]
command=sh -c "php bin/console swarrot:consume:sms_credit_refill_credit sms_credit_refill_credit --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_initial_order]
command=sh -c "php bin/console swarrot:consume:sms_credit_initial_order sms_credit_initial_order --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_refill_order]
command=sh -c "php bin/console swarrot:consume:sms_credit_refill_order sms_credit_refill_order --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sms_credit_alert_consumed_credit]
command=sh -c "php bin/console swarrot:consume:sms_credit_alert_consumed_credit sms_credit_alert_consumed_credit --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0