# infrastructure/services/remote/supervisord/user.conf
[program:user_password]
command=php bin/console swarrot:consume:user_password user_password --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
