# infrastructure/services/remote/supervisord/user_archive.conf
[program:user_archive_request]
command=php bin/console swarrot:consume:user_archive_request user_archive_request --env=prod
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
