# infrastructure/services/remote/supervisord/user.conf
[program:user_password]
command=sh -c "php bin/console swarrot:consume:user_password user_password --env=prod || (sleep 60s && false)"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:user_email]
command=sh -c "php bin/console swarrot:consume:user_email user_email --env=prod"
directory=/var/www
autostart=true
autorestart=true
user=capco
startsecs=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
